# launchd/com.discord-tts-bot.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.discord-tts-bot.bot</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/hyosangahn/Personal-Projects/Discord-TTS-Bot/target/release/tts_bot</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/hyosangahn/Personal-Projects/Discord-TTS-Bot</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>StartInterval</key>
    <integer>60</integer>
    
    <key>StandardOutPath</key>
    <string>/Users/hyosangahn/Personal-Projects/Discord-TTS-Bot/logs/bot.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/hyosangahn/Personal-Projects/Discord-TTS-Bot/logs/bot.error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>
    
    <!-- Wait for database and XTTS to be ready -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>2</integer>
    </dict>
</dict>
</plist>